# supervisord.conf
[supervisord]
nodaemon=true
user=root
loglevel=debug

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[unix_http_server]
file=/var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:mariadb-init]
command=./docker-entrypoint-initdb.d/init-mariadb.sh
autostart=true
autorestart=false
startsecs=0

[program:mariadb-server]
command=/usr/bin/mysqld --user=mysql --console
autostart=false
autorestart=true
user=mysql